# httpd-true-flatpak.socket
[Unit]
Description=HTTPD true

[Socket]
ListenStream=8009

[Install]
WantedBy=default.target
