# httpd-true.socket
[Unit]
Description=HTTPD true

[Socket]
ListenStream=8008

[Install]
WantedBy=default.target
