# httpd-true-flatpak.socket
[Unit]
Description=HTTPD true (flatpak)

[Socket]
ListenStream=8008

[Install]
WantedBy=default.target
