# httpd-true-flatpak.socket
[Unit]
Description=HTTPD true (flatpak)

[Socket]
ListenStream=8009

[Install]
WantedBy=default.target
